# deploy/ghost-server.service
[Unit]
Description=Ghost Protocol Server - Stealth Proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/ghost-server /etc/ghost/ghost-server.conf
Restart=always
RestartSec=5
LimitNOFILE=65535

# Security hardening
NoNewPrivileges=false
ProtectSystem=false
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ghost-server

[Install]
WantedBy=multi-user.target
